# deploy/systemd/pixiv-novel-sync.service
[Unit]
Description=Pixiv Novel Bookmark Sync
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/opt/pixiv-novel-sync/app-redeploy
EnvironmentFile=/opt/pixiv-novel-sync/app-redeploy/.env
ExecStart=/opt/pixiv-novel-sync/app-redeploy/.venv/bin/pixiv-novel-sync --config /opt/pixiv-novel-sync/app-redeploy/config/config.yaml sync-bookmarks
Nice=10

[Install]
WantedBy=multi-user.target
